# linux/roles/wod.k8s-dropcaches/templates/services/k8s-dropcaches.timer
[Unit]
Description=Run Drop Caches Service

[Timer]
OnCalendar=*-*-* 04:10:00
Unit=k8s-dropcaches.service
Persistent=true     
 
[Install]
WantedBy=timers.target