# linux/roles/wod.k8s-dropcaches/templates/services/k8s-dropcaches.service
[Unit]
Description=Drop Caches

[Service]
Type=oneshot
ExecStart=/etc/kubernetes/scripts/k8s-dropcaches.sh